# backups/2025-06-17_v4/scripts/claude-backup-daemon.service
[Unit]
Description=CLAUDE 120-Minute Backup Daemon for Christian
Documentation=file:///Users/scarmatrix/Project/CLAUDE_improvement/CLAUDE.md
After=network.target

[Service]
Type=simple
User=scarmatrix
Group=staff
WorkingDirectory=/Users/scarmatrix/Project/CLAUDE_improvement
ExecStart=/usr/bin/python3 /Users/scarmatrix/Project/CLAUDE_improvement/scripts/backup_daemon.py --start --foreground
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/Users/scarmatrix/Project/CLAUDE_improvement/backup_daemon.log
StandardError=append:/Users/scarmatrix/Project/CLAUDE_improvement/backup_daemon.log

# Environment
Environment=PYTHONPATH=/Users/scarmatrix/Project/CLAUDE_improvement/scripts

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target